# systemd/conduit.service.d/khajubridge.conf
# Drop-in for conduit.service — re-applies KhajuBridge firewall after every Conduit start.
# Install: sudo cp this file /etc/systemd/system/conduit.service.d/khajubridge.conf
#          sudo systemctl daemon-reload
[Service]
ExecStartPost=/bin/bash -c 'sleep 2 && /opt/khajubridge/scripts/apply_firewall.sh'
